# scripts/profile.d/instana_python.sh
#!/bin/bash

export AUTOWRAPT_BOOTSTRAP=instana