# scripts/profile.d/instana_nodejs.sh
#!/bin/bash

export NODE_OPTIONS="${NODE_OPTIONS} --require @instana/collector/src/immediate"